# deploy/com.clinicsiting.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clinicsiting.refresh</string>

  <!-- 每月 1 日 09:00 執行刷新；若當時 Mac 未開機，launchd 會在下次喚醒時補跑 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/chenhungwen/Claude/店面選擇/deploy/refresh.sh</string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/chenhungwen/Claude/店面選擇/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chenhungwen/Claude/店面選擇/logs/launchd.err.log</string>
</dict>
</plist>
